# livegramm.conf
[program:livegrambot]
directory=/home/ubuntu/bots/livegram
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/bots/livegram/logfile_err.log
stdout_logfile=/home/ubuntu/bots/livegram/logfile.log